# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:tor]
command=/usr/bin/tor
